# deploy/asr-stream.conf
; Servicio de transcripción + diarización en STREAMING (directos) — :18006.
;   asr-stream -> 127.0.0.1:18006   (POST /v1/realtime/transcriptions, SSE)
; Corre en diart-venv (CPU, 0 VRAM): diart (pyannote 3.1) hace la diarización
; online y reutiliza el Whisper large-v3 de GPU del backend batch (:18005) para
; el texto. Se expone por el tunnel en asr.llm.botalite.es con path /v1/realtime/*
; (encima del catch-all que va a :18005) y se enruta desde LiteLLM por pass-through.
[program:asr-stream]
command=/bin/bash /workspace/asr/streaming/run-stream.sh
directory=/workspace/asr/streaming
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true
priority=460
stdout_logfile=/workspace/logs/asr-stream.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
redirect_stderr=true
