# deploy/asr.conf
; Ejemplo de programa supervisord para el backend ASR.
; Ajusta las rutas a tu instalación. El servicio escucha solo en 127.0.0.1:18005;
; expónlo a través de un proxy/tunnel con autenticación.
[program:asr]
command=/bin/bash /ruta/a/whisper-asr/run-asr.sh
directory=/ruta/a/whisper-asr
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=40
startretries=3
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
priority=450
stdout_logfile=/var/log/asr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
redirect_stderr=true
